# config/etc/profile.d/bash_aliases.sh
alias ls='ls --color=auto'
alias ll='ls -alih --color=auto'

alias ipkg='nala install'
alias ppkg='nala purge --autoremove'
alias updpkg='nala update'
alias upgpkg='nala upgrade'
